# scripts/launchd/ai.desktop-pilot.bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.desktop-pilot.bridge</string>
  <key>ProgramArguments</key>
  <array>
    <string>__NODE_PATH__</string>
    <string>__INSTALL_DIR__/dist/src/server.js</string>
  </array>
  <key>WorkingDirectory</key>
  <string>__INSTALL_DIR__</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>__HOME__/Library/Logs/desktop-pilot-bridge.out.log</string>
  <key>StandardErrorPath</key>
  <string>__HOME__/Library/Logs/desktop-pilot-bridge.err.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>NODE_ENV</key>
    <string>production</string>
    <key>HOME</key>
    <string>__HOME__</string>
  </dict>
</dict>
</plist>
